# scripts/moss-serve.service.example
[Unit]
Description=Moss AI Agent Serve - LangGraph AI agent backend
After=network.target
Before=moss-qq-bot.service

[Service]
Type=simple
User=YOUR_USERNAME
Environment="PATH=/home/YOUR_USERNAME/.bun/bin:/home/YOUR_USERNAME/.opencode/bin:/home/YOUR_USERNAME/.venv.linux/bin:/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=/PATH/TO/Moss
ExecStart=/home/YOUR_USERNAME/.venv.linux/bin/python -m src.api.server --port 8192
Restart=always
RestartSec=5
StandardOutput=append:/PATH/TO/moss-qq-bot/logs/moss-serve.log
StandardError=append:/PATH/TO/moss-qq-bot/logs/moss-serve-error.log

[Install]
WantedBy=multi-user.target
